# data/g915x-rgb.service
[Unit]
Description=G915 X RGB - Apply keyboard lighting profile
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=oneshot
ExecStartPre=/usr/bin/sleep 2
ExecStart=/usr/bin/python3 -m g915x_rgb.apply
Environment=PYTHONPATH=/home/joshua/Development/g915x-rgb/src
RemainAfterExit=no

# No auto-restart; udev handles reconnection triggers
Restart=no

[Install]
WantedBy=graphical-session.target
