# data/g915x-rgb.service
[Unit]
Description=G915 X RGB - Apply keyboard lighting profile
After=graphical-session.target

[Service]
Type=oneshot
ExecStartPre=/usr/bin/sleep 3
ExecStart=/usr/bin/bash -c 'PYTHONPATH=/home/joshua/Development/g915x-rgb/src /usr/bin/python3 -m g915x_rgb.apply'
Restart=no

[Install]
WantedBy=graphical-session.target
