# data/g915x-rgb.service
[Unit]
Description=G915 X RGB - Apply keyboard lighting profile
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 -m g915x_rgb.apply
Environment=PYTHONPATH=/home/joshua/Development/g915x-rgb/src
RemainAfterExit=no

# Retry on failure (keyboard may not be ready immediately)
Restart=on-failure
RestartSec=3
StartLimitBurst=5
StartLimitIntervalSec=60

[Install]
WantedBy=graphical-session.target
