# data/g915x-rgb.service
[Unit]
Description=G915 X RGB - Keyboard lighting watcher
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/bash -c 'PYTHONPATH=/home/joshua/Development/g915x-rgb/src /usr/bin/python3 -m g915x_rgb.watcher'
Restart=on-failure
RestartSec=10

[Install]
WantedBy=graphical-session.target
